# simpleService.service
[Unit]
Description=Simple Service

[Service]
Type=simple
ExecStart=/etc/simpleService.sh
Restart=on-failure
User=ubuntu

